# com.textpolisher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.textpolisher</string>
    <key>ProgramArguments</key>
    <array>
        <string>__SCRIPT_DIR__/.venv/bin/python</string>
        <string>__SCRIPT_DIR__/text_polisher.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/edwindewit/Library/Logs/text-polisher.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/edwindewit/Library/Logs/text-polisher.log</string>
    <key>WorkingDirectory</key>
    <string>__SCRIPT_DIR__</string>
</dict>
</plist>
